# extras/pyccolo.service
[Unit]
Description=Pyccolo
After=multi-user.target

[Service]
ExecStart=/opt/pyccolo/start.sh
Type=forking

[Install]
WantedBy=graphical.target
